# packaging/systemd/arctern.service
[Unit]
Description=arctern ZFS replication daemon
Documentation=https://github.com/okhsunrog/arctern
After=zfs.target network-online.target
Wants=zfs.target network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/arctern daemon --config /etc/arctern/arctern.toml --socket /run/arctern/arctern.sock
Restart=on-failure
RestartSec=5s

# Needs root for zfs(8)/zpool(8) and to spawn ssh with the replication key.
User=root
Group=root

# SQLite observability state (state.db) and the replication SSH key live here.
StateDirectory=arctern
StateDirectoryMode=0700

# UNIX socket lives under /run; systemd creates it on start, removes on stop.
RuntimeDirectory=arctern
RuntimeDirectoryMode=0755

# Logs to stderr -> journald.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=arctern

# Hardening. ProtectSystem=full leaves /etc readable (config) but
# read-only, which is all the daemon needs; StateDirectory/
# RuntimeDirectory stay writable. /dev/zfs access is unaffected.
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
